# com.garmincoach.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.garmincoach.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rschermer/Documents/garmin-coach/scripts/sync_local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rschermer/Documents/garmin-coach</string>

    <!-- Draai dagelijks 08:15 lokale tijd. launchd haalt een gemiste run in
         zodra de Mac weer wakker is. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/rschermer/Documents/garmin-coach/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rschermer/Documents/garmin-coach/logs/launchd.err.log</string>
</dict>
</plist>
